# systemd/iwset.service
[Unit]
Description=IW set region service

[Service]
ExecStart=/path/to/iwreg.sh

[Install]
WantedBy=multi-user.target
