# systemd/iwset.service
[Unit]
Description=IW set region service

[Service]
ExecStart=/path/to/iwset.sh

[Install]
WantedBy=multi-user.target
