# etc/udev/rules.d/regulatory.rules
# CRDA regulatory rules
KERNEL=="regulatory*", ACTION=="change", SUBSYSTEM=="platform", RUN+="/sbin/crda"
